# repmon.service
[Unit]
Description=Repmon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/repmon

[Install]
WantedBy=multi-user.target
